# docker/cron/crontab
* * * * * cd /var/www && /usr/local/bin/php artisan schedule:run >> /var/www/storage/logs/schedule/schedule.log 2>>/var/www/storage/logs/schedule/schedule.error.log
